# books-api-tech-challenge/supervisord_streamlit.conf
[supervisord]
nodaemon=true

[program:streamlit]
directory=/app/app
command=streamlit run streamlit_app.py --server.address=0.0.0.0 --server.port=8000
autostart=true
autorestart=true
stderr_logfile=/var/log/streamlit.err.log
stdout_logfile=/var/log/streamlit.out.log